# supervisor_deamon.ini
[group:daemon]
programs=daemon-clip-videos, daemon-kafka-consumer

[program:daemon-clip-videos]
directory=/Users/lixin/AIFI/aifi-ml/production/local-video-service
command=/Users/lixin/.pyenv/versions/aifi360/bin/python -m src.daemon_clip_videos
process_name=%(program_name)s-%(process_num)s
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopsignal=INT
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=/Users/lixin/AIFI/aifi-ml/production/local-video-service/logs/supervisor_logs/daemon-clip-videos.log

[program:daemon-kafka-consumer]
directory=/Users/lixin/AIFI/aifi-ml/production/local-video-service
command=/Users/lixin/.pyenv/versions/aifi360/bin/python -m src.daemon_kafka_consumer
process_name=%(program_name)s-%(process_num)s
numprocs=1
autostart=true
autorestart=true
startsecs=5
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=20
stdout_logfile=/Users/lixin/AIFI/aifi-ml/production/local-video-service/logs/supervisor_logs/daemon-kafka-consumer.log
